# cron_data/crontab
# This crontab is reloaded when the container starts up (via dockerfile CMD).

* * * * * echo "hello" >> /mylogs/cron.log 2>&1
* * * * * wget http://127.0.0.1:5000/test_log
